# scripts/startup/startup_script.sh.in
#! /bin/sh
#
### BEGIN INIT INFO
# Provides:          @PROJECT_NAME@
# Required-Start:    $network
# Required-Stop:     $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# X-Interactive:     true
# Short-Description: Start/stop the [@PROJECT_NAME@ program]
### END INIT INFO
#
# Startup script for @PROJECT_NAME@
#
NAME=@PROJECT_NAME@
PATH=/usr/bin
DAEMON=$PATH/$NAME
OPTIONS=
LOCK=/var/lock/${NAME}
PIDSPATH=/var/run/

touch $LOCK*
start)
echo "Starting " $0
start-stop-daemon --start --quiet --background --pidfile $PIDSPATH/$NAME.pid --exec $DAEMON
;;
stop)
echo "Stopping " $0
 #
# Do something to stop the target platform here
#
;;
restart)
echo "Restarting " $0
$0 stop
$0 start
;;
status)
echo "Status Not Implemented"
#
# Do something to process status checks here.
#
;;
*)
echo "Usage: /etc/init.d/$0 {start|stop}"
exit 1
;;
esac
exit 0